# etc/filedump-clean.cron
#!/bin/sh
exec find /var/lib/filedump/files/ -type f -mtime +14 -delete
